# BOX/ansible/roles/apache/templates/httpd.service.j2
.include /lib/systemd/system/httpd.service
[Service]
UMask={{ apache_user_umask }}
